# raspberry/systemd/photo-cleanup-sent.service
[Unit]
Description=Photo Uploader - Limpeza de sessões enviadas do Raspberry
After=multi-user.target

[Service]
Type=oneshot
User=dev
Group=dev
ExecStart=/opt/photo-uploader/bin/cleanup_sent.sh
