# server/systemd/photo-pull-externas.service
[Unit]
Description=Photo Server - Pull EXTERNAS from Raspberry devices
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
User=photo-sync
Group=photo-sync
ExecStart=/opt/photo-server/bin/pull_externas.sh
